# evcc-smartload/rootfs/etc/services.d/smartload/finish
#!/usr/bin/execlineb -S0
# Stop the container if our service exits — matches addon `init: false` semantic.
s6-svscanctl -t /var/run/s6/services
